# server/crontab
0,30 * * * * /bin/sh /minecraft/backup.sh
